# templates/default/doorman.service.erb
[Unit]
Description=Doorman Proxy Auth Service

[Service]
Type=simple
WorkingDirectory=/home/doorman/doorman
ExecStart=/usr/bin/npm start
